# web/debian/postinst
#!/bin/sh
#
# see: dh_installdeb(1)

set -e

case "$1" in
    configure)
        update-rc.d szn-sklik-importweb defaults 30

        PARENT_DIRECTORY=/www/sklik/
        chown -R --reference ${PARENT_DIRECTORY} ${PARENT_DIRECTORY}/importweb/
        chown -R nobody:nogroup /www/sklik/importweb/log
	    chown -R nobody:nogroup /www/sklik/importweb/savedfiles
	    chown -R nobody:nogroup /www/sklik/importweb/tmp
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)

    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
