# web/conf/supervizor-flask.conf
[program:test_python]
command=/usr/bin/python -u /var/www/home.light.control/web/index.py
stdout_logfile=/var/www/home.light.controllogs/stdout_log
redirect_stderr=true
