# web/conf/supervizor-nginx.conf
[program:nginx]
command=/var/www/home.light.control/bin/home.light.control.nginx
stdout_events_enabled=true
stderr_events_enabled=true
